# sparkflows-MCP-Python/uvicorn.service
[Unit]
Description=Sparkflows MCP Python Uvicorn Server
After=network.target

[Service]
Type=simple
User=your-username
WorkingDirectory=/path/to/sparkflows-MCP-Python
Environment="MCP_API_KEY=your-api-key-here"
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 8090
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target





